# vm/nginx-portforward.service
[Unit]
Description=Port-forward nginx-service:80 to localhost:80 for GCP Load Balancer
After=network.target

[Service]
Type=simple
User=ciccio_conti0909
# Consente binding su porta 80 senza eseguire come root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Environment=KUBECONFIG=/home/ciccio_conti0909/.kube/config
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env kubectl port-forward --address=0.0.0.0 svc/nginx-service 80:80 --v=0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
